# redmine/redmine-supervisord.conf
[program:redmine]
environment=RAILS_ENV="production"
command=bundle exec rails s Puma
