# jenkins/jenkins-supervisord.conf
[program:jenkins]
command=java /var/lib/jenkins.war
user=jenkins
