# jenkins/jenkins-supervisord.conf
[program:jenkins]
environment=JENKINS_HOME="/home/jenkins"
command=java -jar /usr/lib/jenkins/jenkins.war
