# jenkins/jenkins-supervisord.conf
[program:jenkins]
user=jenkins
environment=HOME="/home/jenkins",JENKINS_HOME="/home/jenkins"
command=java -jar /usr/lib/jenkins/jenkins.war
