# mysql/mysql-supervisord.conf
[program:mysql]
#user=mysql
command=/usr/bin/mysqld_safe
